# 2018/ros_catkin_ws/devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/pi/RCR/2018/ros_catkin_ws/src/ros/rosbash/rostcsh"
